# etc/cron.d/pragpub-mailer
0 5 * * * root /usr/bin/flock -n /run/pragpub-mailer.lock /usr/bin/pragpub-mailer
